# devops/systemd/aiagent.service
[Unit]
Description=AI Agent Core Service
After=network.target

[Service]
Type=simple
User=aiagent
Group=aiagent
EnvironmentFile=/etc/ai-agent/env
WorkingDirectory=/opt/ai-agent
ExecStart=/opt/ai-agent/venv/bin/gunicorn -k eventlet -w 1 -b 0.0.0.0:8000 dashboard.app:app
Restart=always
RestartSec=5
Environment=ENV=production

[Install]
WantedBy=multi-user.target

